# systemd/certbot-failure-notify.service
[Unit]
Description=Let's Encrypt Certificate Renewal Failure Notification
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=__INSTALL_DIR__/scripts/failure-notify.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target