# systemd/certbot-expiry-check.service
[Unit]
Description=Let's Encrypt Certificate Expiry Check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=__INSTALL_DIR__/scripts/cert-expiry-check.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target